# packages/pyscreamor/pyscreamor-0.0.3.tar.gz/pyscreamor-0.0.3/pyscreamor/pyscreamor.service
[Unit]
Description=screamor service
After=network-online.target
[Service]
User=de
ExecStart=/usr/bin/python3 -m pyscreamor.screamor